# config/etc/supervisor/conf.d/nodejs.conf
[program:nodejs]
directory=/srv/nodejs/upload
command=/usr/bin/node server
stopsignal=QUIT
redirect_stderr=true
autorestart=true
